# cron/crontab
# m h  dom mon dow   command
*/5 * * * * cd /var/www/laravel; /usr/bin/php bin/console app:echo >>/tmp/cronz.log

